# supervisord.conf
[supervisord]
nodaemon=true

[program:conveyor]
command=/usr/bin/conveyor-svc --config /etc/conveyor.conf

[program:makerware]
command=/usr/bin/makerware
